# deploy/roundtable.service
# Install on the server as /etc/systemd/system/roundtable.service
# Then:  sudo systemctl daemon-reload && sudo systemctl enable --now roundtable
#
# First-time server setup (once):
#   sudo mkdir -p /opt/roundtable && sudo chown $USER /opt/roundtable
#   git clone <repo> /opt/roundtable && cd /opt/roundtable
#   python3 -m venv .venv && ./.venv/bin/pip install -r requirements.txt
#   cp .env.example .env   # then edit secrets
#
[Unit]
Description=RoundTable ticket tracker
After=network.target

[Service]
Type=simple
User=roundtable
WorkingDirectory=/opt/roundtable
EnvironmentFile=/opt/roundtable/.env
ExecStart=/opt/roundtable/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
